# images/iso/overlay/etc/systemd/system/getty@tty1.service.d/override.conf.template
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noissue --autologin __DEFAULT_USER__ %I $TERM
Type=idle
