# packaging/deb/crateos-policy/lib/systemd/system/crateos-policy.timer
[Unit]
Description=Run CrateOS policy check periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
