# packaging/deb/crateos-agent/lib/systemd/system/crateos-agent.service
[Unit]
Description=CrateOS Agent — platform state enforcer
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0
OnFailure=crateos-agent-watchdog.service

[Service]
Type=simple
ExecStart=/usr/local/bin/crateos-agent
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crateos-agent

[Install]
WantedBy=multi-user.target
