# packaging/deb/crateos/etc/update-motd.d/00-crateos-header
#!/bin/sh
printf 'CrateOS\n'
printf 'Ubuntu-derived framework appliance\n'
