# packaging/deb/crateos/etc/update-motd.d/10-help-text
#!/bin/sh
printf 'Primary interface: CrateOS control surface\n'
printf 'Base platform: Ubuntu noble derivative\n'
printf 'Admin shell access is break-glass only.\n'
